# db/dist/etc/supervisor/conf.d/init.db.conf
[program:dbinit]
command=/usr/local/sbin/init.db.sh
autostart=true
autorestart=false
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
