# db/dist/etc/supervisor/conf.d/mysql.conf
[program:mysql]
command=/usr/sbin/mysqld
autostart=true
autorestart=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
